# linux/systemd/timer/test_change_script.service
[Unit]
Description=test change script

[Service]
Type=simple
ExecStart=/usr/bin/print_hello.sh
